# etc/sv/sshd/run
#!/usr/bin/env bash
set -euo pipefail

echo "Starting SSH Daemon..."
exec 2>&1
exec /usr/sbin/sshd -D -e
