# etc/sv/plex/run
#!/usr/bin/env bash
set -euo pipefail

source /etc/environment

chown plex:plex -R "${PLEX_HOME}"

if [[ ! $(ps -A | grep avahi-daemon) ]]; then
    echo 'Avahi has not started, not starting Plex yet.'
    exit 1
fi

echo "Starting Plex Media Server..."
exec 2>&1
exec chpst -u plex /usr/sbin/start_pms
